# files/etc/init.d/overlay-fix
#!/bin/sh /etc/rc.common
START=15

start() {
    mount | grep -q "overlay" && exit 0
    /sbin/block mount
}
